# external/packages/blobs/desktop/skel_new/.config/autostart/org.bananapi.bpiwiz.desktop
[Desktop Entry]
Name=opiwiz
Exec=sudo opiwiz
Icon=org.bananapi.opiwiz
Terminal=false
Type=Application
Categories=GTK;
StartupNotify=true
