# external/packages/bsp/common/etc/cron.d/bananapi-updates
@reboot root /usr/lib/bananapi/bananapi-apt-updates
@daily root /usr/lib/bananapi/bananapi-apt-updates
